# deploy/gloski.service
[Unit]
Description=Gloski Server Control Dashboard
After=network.target

[Service]
Type=simple
User=your-username
ExecStart=%h/.local/bin/gloski -c %h/.config/gloski/config.json
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
